# rootfs/etc/systemd/system/systemd-growfs@sysroot.service.d/zfs.service
[Unit]
Description=Load ZFS module
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/sbin/modprobe zfs

[Install]
WantedBy=multi-user.target