# rootfs/common/etc/profile.d/systemctl_helpers.sh
alias sctl='run0 systemctl'
alias uctl='systemctl --user'