# rootfs/common/usr/lib/systemd/system/bootc-fetch-apply-updates.timer
[Unit]
Description=Apply bootc updates
Documentation=man:bootc(8)
ConditionPathExists=/run/ostree-booted

[Timer]
# This time is relatively arbitrary and obviously expected to be overridden/changed
OnUnitInactiveSec=72h
# When deploying a large number of systems, it may be beneficial to increase this value to help with load on the registry.
RandomizedDelaySec=30m
OnCalendar=*-*-* 06:00:00

[Install]
WantedBy=timers.target